# templates/supervisor.rpirest_tasks.conf.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

[program: rpirest_async_tasks]
command=/usr/local/bin/celery -A rpirest worker -l info
directory={{ rpirest_app_path }}
user={{ rpirest_app_user }}
autostart=true
autorestart=true
redirect_stderr=True
