# templates/supervisor.pirest_sense_hat.conf.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

[program:{{app_name}}]
command=/usr/local/bin/gunicorn {{app_name}}.wsgi --workers=4 -t 300 --bind 127.0.0.1:{{ app_port }}
directory={{ app_path }}
user={{ app_user }}
autostart=true
autorestart=true
redirect_stderr=True
