# templates/supervisor.pirest_sense_hat.conf.j2
# This file was generated by Ansible for {{ansible_fqdn}}
# Do NOT modify this file by hand!

[program:{{rpirest_app_name}}]
command=/usr/local/bin/gunicorn {{rpirest_app_name}}.wsgi --workers=4 -t 300 --bind 127.0.0.1:{{ rpirest_app_port }}
directory={{ rpirest_app_path }}
user={{ rpirest_app_user }}
autostart=true
autorestart=true
redirect_stderr=True
